# packages/keeper/deploy/seawall-keeper-fuel.timer
[Unit]
Description=Run the Seawall gas auto-faucet every 30 min

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
